# fan-control.service
[Unit]

Description=Fan-Control

[Service]

ExecStart=/usr/sbin/fan-control
